# packages/borgee-helper/install/cloud.borgee.host-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>cloud.borgee.host-bridge</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/sandbox-exec</string>
      <string>-f</string>
      <string>/Library/Application Support/Borgee/borgee-helper.sb</string>
      <string>/usr/local/bin/borgee-helper</string>
      <string>--socket=/Users/Shared/Borgee/borgee-helper.sock</string>
      <string>--audit-log=/Library/Logs/Borgee/audit.log.jsonl</string>
      <string>--grants-db=file:/Library/Application Support/Borgee/server.db?mode=ro&amp;_busy_timeout=5000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>UserName</key>
    <string>_borgee-helper</string>

    <key>GroupName</key>
    <string>_borgee-helper</string>

    <key>StandardOutPath</key>
    <string>/Library/Logs/Borgee/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/Borgee/stderr.log</string>
  </dict>
</plist>
